# freebsd/pkg-plist
etc/fidod.conf.sample
etc/rc.d/fidod.sh.sample
sbin/fidod
share/doc/fidod/NEWS
share/doc/fidod/README
share/doc/fidod/TODO
info/fidod.info
